# test/fixtures/cron-agent-capped/agent.cron
@daily timeout 900 claude -p "nightly triage" --max-turns 5
0 6 * * 1 /usr/local/bin/backup.sh
