# test/fixtures/cron-agent-capped/deploy.crontab
*/15 * * * * claude -p "check the queue" --max-turns 2
