# test/fixtures/loop-indirect-npm/crontab
*/15 * * * * cd /srv/repo && npm run agent:fix >> /var/log/agent.log 2>&1
