# test/fixtures/cron-no-agent/crontab
# Ordinary crontab — no agent anywhere. Must stay N/A.
0 3 * * * /usr/local/bin/backup.sh --verbose
@reboot /usr/local/bin/warm-cache.sh
