# test/fixtures/loop-systemd/agent.service
[Service]
Type=oneshot
ExecStart=/usr/bin/claude -p "triage the queue"
