# test/fixtures/loop-systemd/review.service
[Service]
RuntimeMaxSec=900
ExecStart=/usr/bin/claude -p "review the open PRs"
