# test/fixtures/loop-systemd/review.timer
[Timer]
OnCalendar=hourly
Unit=review.service
